# systemd/urbit-meld@.timer
[Unit]
Description=Run urbit-meld@%i weekly

[Timer]
OnCalendar=Mon *-*-* 07:30:00

[Install]
WantedBy=timers.target
